# upstart/cmake_stubs/pgbouncer.conf.in
# PgBouncer Upstart script

description "pgbouncer"

start on started postgresql

stop on runlevel [016]

respawn

respawn limit 10 30

expect fork

script
    exec pgbouncer -u @UNIX_USER@ @concordia-server_SOURCE_DIR@/db/pgbouncer.ini &
    sleep 6
end script
